# supervisord.conf
[supervisord]
nodaemon=true

[program:pypiserver]
command=pypi-server -p 8080 /packages
autorestart=false
exitcodes=0

[program:sync_to_oss]
command=/sync_to_oss.sh
autorestart=false
exitcodes=0
