# supervisord.conf
[supervisord]
nodaemon=true

[program:pypiserver]
command=/start_pypiserver.sh
autorestart=false
exitcodes=0

[program:sync_to_oss]
command=/sync_to_oss.sh
autorestart=false
exitcodes=0
